# templates/systemd/systemd-journald-autovacuum.service.j2
[Unit]
Description=Run journalctl vacuuum procedure

[Service]
Type=oneshot
ExecStart=/bin/journalctl --vacuum-time={{ system_cleanup_logs_retention }}
